# 02_recursos/gestionalmacen.socket
[Unit]
Description=Socket de la API de Gestión Almacén

[Socket]
ListenStream=/run/gestionalmacen.sock

[Install]
WantedBy=sockets.target
